# systemd/user/db0sue.service
[Service]
Type=exec
WorkingDirectory=/cb/logbook.rrdxa.org/rrdxa-log
ExecStart=/cb/logbook.rrdxa.org/rrdxa-log/gather/cluster db0sue.de
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
