# systemd/user/dxsummit.service
[Service]
Type=exec
WorkingDirectory=/cb/logbook.rrdxa.org/rrdxa-log
ExecStart=/cb/logbook.rrdxa.org/rrdxa-log/gather/dxsummit
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
