# systemd/user/pskreporter.service
[Service]
Type=exec
WorkingDirectory=/cb/logbook.rrdxa.org/rrdxa-log
ExecStart=/cb/logbook.rrdxa.org/rrdxa-log/gather/pskreporter
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
