# systemd/user/qo100ft.service
[Service]
Type=exec
WorkingDirectory=/cb/logbook.rrdxa.org/rrdxa-log
ExecStart=/cb/logbook.rrdxa.org/rrdxa-log/gather/qo100ft
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
