# systemd/user/telnetsrv.socket
[Socket]
ListenStream=8000
Accept=true

[Install]
WantedBy=sockets.target
